# startup/lxc/openrc-init/printmsg
#!/sbin/openrc-run

description="Prints a msg"

start() {
	printf "bbbooted"
}